# BCKPDZ/ansible/roles/client/templates/borg-backup.service
[Unit]
Description=Borg Backup of /etc

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borg-backup.sh
